# config/supervisord.conf
[supervisord]
nodaemon=true

[program:znc]
command=sudo -u znc -g znc /usr/bin/znc --datadir /var/lib/znc
autorestart=true
